# roles/atomix-cluster/templates/atomix.service.j2
[Unit]
Description=Atomix Cluter

[Service]
ExecStart=/opt/atomix/bin/atomix-agent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
